# mongodb.sv.conf
[program:mongodb]
command=/bin/bash -c "service mongodb restart"
redirect_stderr=true